# debian/prerm
systemctl stop crowdsec-cloudflare-bouncer || echo "cannot stop service"
systemctl disable crowdsec-cloudflare-bouncer || echo "cannot disable service"
